# snippets/2022/X728/OTA/x728rootfs/etc/systemd/system/munin-node@.service
[Unit]
Description = munin-node service

[Service]
ExecStart = /bin/munin-node

StandardInput = socket

