# snippets/2022/X728/OTA/x728rootfs/etc/systemd/system/x728ups.service
[Unit]
Description=X728 UPS Monitor

[Service]
ExecStart=/etc/x728/upsmon.sh

[Install]
WantedBy=multi-user.target
