# snippets/X728/OTA/x728rootfs/etc/systemd/system/munin-node.socket
[Unit]
Description = munin-node

[Socket]
ListenStream = 4949
Accept = yes

[Install]
WantedBy = sockets.target

