# snippets/xwin-hotplug/99-xwin-hotplug.rules
ACTION=="change", KERNEL=="card*", SUBSYSTEM=="drm", RUN+="/usr/local/bin/xwin-hotplug"

